# examples/example4/caddy.container
[Container]
Image=docker.io/library/caddy
Network=example4-net.network
ContainerName=caddy-container
[Install]
WantedBy=default.target
