# examples/example4/example4.socket
[Unit]
Description=Example 4 socket

[Socket]
ListenStream=0.0.0.0:80

[Install]
WantedBy=sockets.target
