# examples/example7/whoami1.container
[Container]
Image=docker.io/traefik/whoami
ContainerName=whoami1
Network=example7.network
