# examples/example2/example2.socket
[Unit]
Description=Example 2

[Socket]
ListenStream=0.0.0.0:80

[Install]
WantedBy=sockets.target
