# examples/example6/example6-proxy.container.in
[Unit]
Requires=user@${envsubst_uid}.service
After=user@${envsubst_uid}.service
RequiresMountsFor=/run/user/${envsubst_uid}/containers
Requires=example6-backend.socket
After=example6-backend.socket

[Service]
User=${envsubst_user}

[Container]
Environment=NGINX=3;
Image=docker.io/library/nginx
Network=none
SecurityLabelDisable=true
User=0:0
UserNS=keep-id:uid=101,gid=101
Volume=/home/${envsubst_user}/nginx-reverse-proxy-conf:/etc/nginx/conf.d:Z
Volume=/run/user/${envsubst_uid}/backend-socket:/var/socket:Z
